# packages/nginx-amplify-agent/deb/debian/nginx-amplify-agent.conffiles
/etc/init.d/amplify-agent
/etc/logrotate.d/amplify-agent
